# test/unit/test_service/testa.service
[Unit]
Description=DBusManager Test Service A

[Service]
Type=simple
ExecStart=/usr/local/bin/testa
Restart=no

[Install]
WantedBy=multi-user.target
